# system/etc/systemd/system/user@.service.d/10-oomd.conf
# Monitor the whole user session for sustained memory-pressure and kill the
# heaviest descendant cgroup (a container, firefox, or a runaway AI pane) —
# NOT the compositor, which session.slice marks ManagedOOMPreference=avoid.
# Mirrors Fedora's default user@.service oomd activation.
[Service]
ManagedOOMMemoryPressure=kill
