# system/etc/systemd/system/vector.service.d/20-no-net-block.conf
# Soften vector's network gate. The packaged unit pulls network-online.target
# onto the boot critical-chain via a hard `Requires=`. Vector retries its Loki
# sink on its own, so a soft `Wants=` (with the inherited ordering-only After=)
# is correct and keeps vector off the critical path. Vector is also disabled at
# boot now (on-demand via the Mod+5 monitoring scene), so this is belt-and-
# suspenders for the on-demand start path. Resetting Requires= is safe:
# network-online.target was its only hard requirement.
[Unit]
Requires=
Wants=network-online.target
