# system/etc/systemd/system/grafana-alloy.service.d/10-ebpf-caps.conf
# [ENGINE] monitoring — enable pyroscope.ebpf continuous profiling.
# Alloy runs as the unprivileged `grafana-alloy` user. eBPF profiling needs:
#   CAP_BPF             load BPF programs/maps (unprivileged_bpf_disabled stays on)
#   CAP_PERFMON         perf_event_open for sampling (bypasses perf_event_paranoid)
#   CAP_SYS_PTRACE      read /proc/<pid>/exe|maps|root for discovery + symbolization
#   CAP_SYSLOG          read /proc/kallsyms addresses (needs kptr_restrict=1)
#   CAP_DAC_READ_SEARCH read tracepoint IDs under /sys/kernel/tracing (root mode 700)
# RLIMIT_MEMLOCK is raised by systemd (no CAP_SYS_RESOURCE needed). This is the
# least-privilege set in lieu of root / CAP_SYS_ADMIN; CapabilityBoundingSet is
# tightened to exactly these. Pairs with /etc/sysctl.d/90-pyroscope-ebpf.conf.
[Service]
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_SYSLOG CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_SYSLOG CAP_DAC_READ_SEARCH
LimitMEMLOCK=infinity
