# system/etc/systemd/system/systemd-networkd-wait-online.service.d/10-fast.conf
# Boot perf vs. correctness: `--any` returns as soon as ONE required link is
# routable (only enp0s31f6 is RequiredForOnline; wlan/wwan are =no). The old
# `--timeout=5` cap was SHORTER than the real ~10.6s ethernet DHCP, so the unit
# FAILED every boot and marked the system `degraded` (the orange [FAILED]).
# Fix: 20-ethernet.network now sets RequiredFamilyForOnline=ipv4 (don't wait on
# IPv6 to call it online), and we raise the cap to 20s so a slightly-slow DHCP
# still SUCCEEDS instead of failing. `--any` means we still return the instant
# IPv4 is routable, so the typical boot cost is unchanged. wait-online stays
# AVAILABLE (not masked) so wg-quick@wg0 / tailscale still get an online signal.
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any --timeout=20
