# system/etc/systemd/system/systemd-networkd-wait-online.service.d/10-fast.conf
# Boot perf: don't block boot on full network readiness.
# `--any` returns as soon as ONE managed link is routable; `--timeout=5` caps
# the wait at 5s (was an unbounded ~10.6s DHCP wait on enp0s31f6 — the top of
# the boot critical-chain). wait-online stays AVAILABLE (deliberately not masked)
# so units that legitimately gate on connectivity (wg-quick@wg0, tailscale,
# future) still get an online signal.
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any --timeout=5
