# systemd-user/.config/systemd/user/discord-mcp.service
[Unit]
Description=Discord MCP server (SL1C3D-L4BS bot wrapper)
Documentation=file://%h/SlicedLabs/tools/discord/README.md
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/uv run --quiet --script %h/SlicedLabs/tools/discord/server.py
WorkingDirectory=%h/discord-mcp
Restart=on-failure
RestartSec=5s
# Token comes from `pass` at startup — no EnvironmentFile, no .env.

# Hardening: this process only needs the network and the user's pass+gpg
# state. Block everything else.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/discord-mcp %h/.password-store %h/.gnupg
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Not enabled by default — see README runbook. Claude Code MCP harness
# may already launch this; verify before `systemctl --user enable --now`.

[Install]
WantedBy=default.target
