# systemd-user/.config/systemd/user/sl-alert-health.timer
[Unit]
Description=Run the SlicedLabs health watchdog periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
