# systemd-user/.config/systemd/user/sl-alert-health.service
[Unit]
Description=SlicedLabs health watchdog — failed units · disk → sl-alert
After=ntfy.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sl-alert-health
